# scripts/com.jhizzard.termdeck-supervise.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  TermDeck stack supervisor — runs scripts/termdeck-supervise.sh every 60s and at
  login, keeping the server / Mnestra webhook / cloudflared tunnel / MCP bridge
  alive. Install (operator, once — launchd cannot be loaded from the sandbox):

    cp scripts/com.jhizzard.termdeck-supervise.plist ~/Library/LaunchAgents/
    launchctl load -w ~/Library/LaunchAgents/com.jhizzard.termdeck-supervise.plist
    # to stop:  launchctl unload -w ~/Library/LaunchAgents/com.jhizzard.termdeck-supervise.plist

  Machine-specific: the absolute paths below assume the repo at the default
  location and HOME=/Users/joshuaizzard. Edit both if either differs.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jhizzard.termdeck-supervise</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/joshuaizzard/Documents/Graciella/ChopinNashville/SideHustles/TermDeck/termdeck/scripts/termdeck-supervise.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>60</integer>
  <key>StandardOutPath</key>
  <string>/Users/joshuaizzard/.termdeck/logs/supervise.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/joshuaizzard/.termdeck/logs/supervise.err.log</string>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
